# vitals/scripts/gunicorn.service
[Unit]
Description=Gunicorn for healthstack (Django)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/srv/healthstack
Environment="PYTHONUNBUFFERED=1"
EnvironmentFile=-/etc/default/healthstack
ExecStart=/srv/healthstack/.venv/bin/gunicorn \
  --pythonpath /srv/healthstack \
  --workers 3 \
  --threads 4 \
  --timeout 120 \
  --bind 127.0.0.1:8000 \
  healthstack.wsgi:application
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
